# debian/reminder.cron.d
* * * * * root if /usr/sbin/ngcp-check_active -q; then /usr/sbin/reminder 1>/dev/null; fi
